# systemd/targetclid.service
[Unit]
Description=Targetcli daemon
Documentation=man:targetcli(8)
After=network.target targetcli.service

[Service]
Type=simple
ExecStart=/usr/bin/targetclid
Restart=on-failure

[Install]
WantedBy=multi-user.target
Also=targetclid.socket
