# systemd/targetclid.socket
[Unit]
Description=targetclid socket
PartOf=targetclid.service

[Socket]
ListenStream=/var/run/targetclid.sock

[Install]
WantedBy=sockets.target
